# supervisord.conf
[program:foo]
directory=/home/chris/foo
command=/home/chris/.envs/foo/bin/python manage.py runserver "0.0.0.0:8000"